# src/sample/sample.service
[Unit]
Description=YaizuSample!!
After=network.service

[Service]
Type=simple
Restart=on-success
ExecStart=/usr/bin/sample

[Install]
WantedBy=multi-user.target

